# services/redclaw-zeroclaw.service
[Unit]
Description=RedClaw Assistant Runtime
Documentation=https://github.com/nebullii/redclaw
Wants=network-online.target
After=network-online.target redclaw-model.service
Requires=redclaw-model.service

[Service]
Type=simple
User=redclaw
Group=redclaw
EnvironmentFile=-/etc/redclaw/redclaw-zeroclaw.env
Environment=ZEROCLAW_CONFIG_DIR=/etc/redclaw
Environment=ZEROCLAW_WORKSPACE=/var/lib/redclaw/workspace
WorkingDirectory=/var/lib/redclaw
StateDirectory=redclaw
LogsDirectory=redclaw
RuntimeDirectory=redclaw
ExecStart=/usr/local/bin/zeroclaw daemon
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=/var/lib/redclaw /var/log/redclaw /etc/redclaw

[Install]
WantedBy=multi-user.target
