# services/vllm.service
[Unit]
Description=redclaw vLLM Service
After=network.target

[Service]
Type=simple
User=redclaw
Group=redclaw
WorkingDirectory=/opt/redclaw
ExecStart=/usr/local/bin/vllm serve Qwen/Qwen2.5-7B-Instruct --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/redclaw

[Install]
WantedBy=multi-user.target
